# deployment/quadlet-templates/evi-pgadmin.container
# Version: 1.1.0
# Purpose: Podman Quadlet container unit for pgAdmin 4.
# Deployment file: evi-pgadmin.container
# Logic: Runs pgAdmin web interface, connected to internal network.
#        Enabled only if EVI_PGADMIN_ENABLED=true.
#        Accessible ONLY from localhost (127.0.0.1) for security.
#
# Changes in v1.1.0:
# - Port bound to localhost only (127.0.0.1) for security
# - Disabled pgAdmin authentication (SERVER_MODE=False) - safe for localhost-only access
# - Removed PGADMIN_DEFAULT_EMAIL and PGADMIN_DEFAULT_PASSWORD
# - Added servers.json mount for pre-configured evi-db connection
#
# Changes in v1.0.0:
# - Initial pgAdmin container unit
#

[Unit]
Description=evi pgAdmin 4 (Podman)
Wants=network-online.target
After=network-online.target
Requires=evi-network.service evi-db.service
After=evi-network.service evi-db.service

[Container]
# Rendered by evictl
Image={{EVI_PGADMIN_IMAGE}}
ContainerName=evi-pgadmin

# Internal network
Network={{EVI_NETWORK}}

# Publish port to localhost ONLY (security: no external access)
PublishPort={{EVI_PGADMIN_HOST}}:{{EVI_PGADMIN_PORT}}:80

# Environment - Desktop mode (no login required, safe for localhost-only)
Environment=PGADMIN_CONFIG_SERVER_MODE=False
Environment=PGADMIN_CONFIG_MASTER_PASSWORD_REQUIRED=False

# Mount pre-configured servers.json for evi-db connection
Volume={{EVI_STATE_DIR}}/pgadmin/servers.json:/pgadmin4/servers.json:ro,Z

# Resource limits (optional - keep minimal)
PodmanArgs=--memory=512m

[Service]
Restart=always
TimeoutStartSec=180

[Install]
WantedBy=default.target
