# deployment/quadlet-templates/evi-reverse-proxy.container
# Version: 1.2.0
# Purpose: Podman Quadlet container unit for evi reverse proxy (Caddy).
# Deployment file: evi-reverse-proxy.container
# Logic: Terminates TLS on host 80/443 and routes / -> evi-fe, /api -> evi-be.
#
# Changes in v1.2.0:
# - Renamed from evi-proxy to evi-reverse-proxy for consistency.
#
# Changes in v1.1.0:
# - Renamed from evi-reverse-proxy to evi-proxy for consistency.
#
# Changes in v1.0.0:
# - Initial Caddy proxy container (single ingress point)
#

[Unit]
Description=evi reverse proxy (Caddy) (Podman)
Wants=network-online.target
After=network-online.target
Requires=evi-network.service
After=evi-network.service

[Container]
# Rendered by evictl
Image={{EVI_PROXY_IMAGE}}
ContainerName=evi-reverse-proxy

# Internal network (so it can reach evi-fe and evi-be by container name)
Network={{EVI_NETWORK}}

# Publish host ports (Caddy is the only container publishing these)
PublishPort={{EVI_HTTP_PORT}}:80
PublishPort={{EVI_HTTPS_PORT}}:443

# Caddy configuration and state
Volume={{EVI_STATE_DIR}}/proxy/Caddyfile:/etc/caddy/Caddyfile:ro,Z
Volume=evi_caddy_data:/data:Z
Volume=evi_caddy_config:/config:Z

# Optional: manual TLS cert/key mounted as read-only files when enabled by evictl
{{EVI_PROXY_TLS_MOUNTS}}

# Resource limits (optional)
{{EVI_PODMANARGS_LIMITS_PROXY}}

[Service]
Restart=always
TimeoutStartSec=180

[Install]
WantedBy=default.target