# deploy/quadlet-templates/evi-db.container
# Version: 1.0.2
# Purpose: Podman Quadlet container unit for evi PostgreSQL database.
# Deployment file: evi-db.container
# Logic: Runs DB on the internal network only (no host port published).
#
# Changes in v1.0.1:
# - Fixed healthcheck to use explicit postgres user (fixes "role root does not exist" errors)
# - Added EVI_ADMIN_DB_USERNAME environment variable for configurable admin username
#
# Changes in v1.0.0:
# - Initial DB container unit (internal-only, persistent volume)
#
# Changes in v1.0.2:
# - Volume mount changed from /var/lib/postgresql/data to /var/lib/postgresql for PG 18+ compatibility
#

[Unit]
Description=evi database (Podman)
Wants=network-online.target
After=network-online.target
Requires=evi-network.service evi-db-volume.service
After=evi-network.service evi-db-volume.service

[Container]
# Rendered by install.sh (deploy)
Image={{EVI_DB_IMAGE}}
ContainerName=evi-db

# Internal network
Network={{EVI_NETWORK}}

# Persistent data volume (PG 18+ uses /var/lib/postgresql/18/main internally)
Volume=evi_db_data:/var/lib/postgresql:Z

# Environment
Environment=POSTGRES_PASSWORD={{EVI_POSTGRES_PASSWORD}}
Environment=POSTGRES_DB={{EVI_POSTGRES_DB}}
Environment=EVI_APP_DB_PASSWORD={{EVI_APP_DB_PASSWORD}}
Environment=EVI_ADMIN_DB_USERNAME={{EVI_ADMIN_DB_USERNAME}}
Environment=EVI_ADMIN_DB_PASSWORD={{EVI_ADMIN_DB_PASSWORD}}
Environment=SEED_DEMO_DATA={{EVI_SEED_DEMO_DATA}}

# Healthcheck - use explicit postgres user to avoid "role root does not exist" errors
PodmanArgs=--health-cmd="pg_isready -U postgres -d maindb"
PodmanArgs=--health-interval=10s
PodmanArgs=--health-timeout=5s
PodmanArgs=--health-retries=10

# Resource limits (optional)
{{EVI_PODMANARGS_LIMITS_DB}}

[Service]
Restart=always
TimeoutStartSec=180

[Install]
WantedBy=default.target


