# deploy/quadlet-templates/evi-be.container
# Version: 1.1.0
# Purpose: Podman Quadlet container unit for evi backend API.
# Deployment file: evi-be.container
# Logic: Runs backend only on internal network; exposed externally via evi-reverse-proxy (Caddy) as /api over HTTPS.
#
# Changes in v1.1.0:
# - Backend starts after DB (After=evi-db.service, Wants=evi-db.service) so container is created only when DB exists
#
# Changes in v1.0.0:
# - Initial backend container unit (no host port, JWT key mounted via podman secret)
#

[Unit]
Description=evi backend (Podman)
Wants=network-online.target
After=network-online.target
Requires=evi-network.service
After=evi-network.service
Wants=evi-db.service
After=evi-db.service

[Container]
# Rendered by install.sh (deploy)
Image={{EVI_BE_IMAGE}}
ContainerName=evi-be

# Internal network
Network={{EVI_NETWORK}}

# Environment (DB connectivity)
Environment=NODE_ENV={{EVI_NODE_ENV}}
Environment=API_PORT={{EVI_BE_PORT}}
Environment=PGHOST=evi-db
Environment=PGPORT={{EVI_DB_PORT}}
Environment=PGDATABASE={{EVI_POSTGRES_DB}}
Environment=PGUSER=app_service
Environment=PGPASSWORD={{EVI_APP_DB_PASSWORD}}

# CORS (must include https://<domain>)
Environment=EVI_CORS_ORIGINS={{EVI_CORS_ORIGINS}}

# JWT private key from Podman secret mounted as file
PodmanArgs=--secret={{EVI_JWT_SECRET_NAME}},type=mount,target=/run/secrets/{{EVI_JWT_SECRET_NAME}}
Environment=JWT_PRIVATE_KEY_PATH=/run/secrets/{{EVI_JWT_SECRET_NAME}}

# Resource limits (optional)
{{EVI_PODMANARGS_LIMITS_BE}}

[Service]
Restart=always
TimeoutStartSec=180

[Install]
WantedBy=default.target


