# deploy/quadlet-templates/evi-reverse-proxy.container
# Version: 1.3.0
# Purpose: Podman Quadlet container unit for evi reverse proxy (Caddy).
# Deployment file: evi-reverse-proxy.container
# Logic: Terminates TLS on host 80/443 and routes / -> evi-fe, /api -> evi-be.
#
# Changes in v1.3.0:
# - Fixed network binding: replaced PublishPort with PodmanArgs --publish 0.0.0.0:PORT:PORT
# - Ensures ports listen on all interfaces (0.0.0.0) for external access
# - Makes configuration explicit and reliable across Podman versions
#
# Changes in v1.2.0:
# - Renamed from evi-proxy to evi-reverse-proxy for consistency.
#
# Changes in v1.1.0:
# - Renamed from evi-reverse-proxy to evi-proxy for consistency.
#
# Changes in v1.0.0:
# - Initial Caddy proxy container (single ingress point)
#

[Unit]
Description=evi reverse proxy (Caddy) (Podman)
Wants=network-online.target
After=network-online.target
Requires=evi-network.service
After=evi-network.service

[Container]
# Rendered by install.sh (deploy)
Image={{EVI_PROXY_IMAGE}}
ContainerName=evi-reverse-proxy

# Internal network (so it can reach evi-fe and evi-be by container name)
Network={{EVI_NETWORK}}

# Publish host ports on all interfaces (0.0.0.0) for external access
# Note: PublishPort in Quadlet doesn't support IP binding, using PodmanArgs instead
PodmanArgs=--publish 0.0.0.0:{{EVI_HTTP_PORT}}:80
PodmanArgs=--publish 0.0.0.0:{{EVI_HTTPS_PORT}}:443

# Caddy configuration and state
Volume={{EVI_STATE_DIR}}/reverse-proxy/Caddyfile:/etc/caddy/Caddyfile:ro,Z
Volume=evi_caddy_data:/data:Z
Volume=evi_caddy_config:/config:Z

# Optional: manual TLS cert/key mounted as read-only files when enabled by install.sh (deploy)
{{EVI_PROXY_TLS_MOUNTS}}

# Resource limits (optional)
{{EVI_PODMANARGS_LIMITS_PROXY}}

[Service]
Restart=always
TimeoutStartSec=180

[Install]
WantedBy=default.target